# e/easyvdr-3.99.00/etc/systemd/easyvdr/easyvdr-create-ramdisk.service
[Unit]
Description=Easyvdr Create Ramdisk
After=easyvdr-vdr.service 

[Service]
Type=simple
#ExecStartPre=
ExecStart=/usr/lib/easyvdr/easyvdr-start/systemd/easyvdr-create-ramdisk
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
