# e/easyvdr-3.99.00/etc/systemd/easyvdr/easyvdr-save-ramdisk.service
[Unit]
Description=Easyvdr save Ramdisk

[Service]
Type=simple
ExecStart=/usr/lib/easyvdr/easyvdr-start/systemd/easyvdr-safe-ramdisk
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=exit.target