# e/easyvdr-3.99.00/startvdr
#! /bin/bash

### BEGIN INIT INFO
# Provides:          startvdr
# Short-Description: EasyVDR-script to start vdr for non root-user
# Description:       start vdr
### END INIT INFO

# Author: Uwe Kiehl <uwe@easy-vdr.de>
# adopted to easyvdr 2.0 by michel8

# Do NOT "set -e"

#set -x

. /usr/share/easyvdr/setup/easyvdr-setup-defines

#####
Get_Frontend_Class()
{
case $VDR_FRONTEND in
     "vdr-sxfe") FRONTEND_CLASS="USE_X11" ;;
         "xine") FRONTEND_CLASS="USE_X11" ;;
         "ENNA") FRONTEND_CLASS="USE_X11" ;;
         "KODI") FRONTEND_CLASS="USE_X11" ;;
 "Softhddevice") FRONTEND_CLASS="USE_X11" ;;
       "DVB-FF") FRONTEND_CLASS="NON_X11" ;;
     "DVBS2-FF") FRONTEND_CLASS="NON_X11" ;;
      "reelbox") FRONTEND_CLASS="NON_X11" ;;
       "PVR350") FRONTEND_CLASS="NON_X11" ;;
       "keines") FRONTEND_CLASS="NON_FRE" ;;
              *) FRONTEND_CLASS="NON_FRE" ;;
esac
}

Get_Frontend_Class
if [ "x$FRONTEND_CLASS" == "xUSE_X11" ]; then
    chvt 7
fi

sudo /usr/lib/easyvdr/easyvdr-start/easyvdr-vdr start
sudo /usr/lib/easyvdr/easyvdr-start/easyvdr-frontend start
