# e/easyvdr-3.99.00/debian/postrm
#! /bin/sh
# 
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postrm> `remove'
#        * <postrm> `purge'
#        * <old-postrm> `upgrade' <new-version>
#        * <new-postrm> `failed-upgrade' <old-version>
#        * <new-postrm> `abort-install'
#        * <new-postrm> `abort-install' <old-version>
#        * <new-postrm> `abort-upgrade' <old-version>
#        * <disappearer's-postrm> `disappear' <r>overwrit>r> <new-version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package



remove_service()
{
 #Ab hier easyvdr service
 SERVICE[0]="/etc/systemd/easyvdr/easyvdr-systemstart.service"
 SERVICE[1]="/etc/systemd/easyvdr/easyvdr-alsa-restore.service"
 SERVICE[2]="/etc/systemd/easyvdr/easyvdr-cleanscreen.service"
 SERVICE[3]="/etc/systemd/easyvdr/easyvdr-create-ramdisk.service"
 SERVICE[4]="/etc/systemd/easyvdr/easyvdr-save-ramdisk.service"


 for CHK_SERVICE in ${SERVICE[*]}
 do
  #Pfad vom Service Namen abschneiden
  SERVICE_NAME=${CHK_SERVICE##*/}

  if [ -f $CHK_SERVICE ]; then
   STAT=$(systemctl is-active "$SERVICE_NAME") || true > /dev/null 2>&1
   if [ $STAT != "inactive" ]; then
    systemctl stop "$SERVICE_NAME" || true > /dev/null 2>&1
   fi

   if [ -f "$CHK_SERVICE" ] ;then
    STAT=$(systemctl is-enabled "$SERVICE_NAME") || true > /dev/null 2>&1
    if [ $STAT == "enabled" ]; then
     systemctl disable "$SERVICE_NAME" || true > /dev/null 2>&1
    fi
   fi
  fi
 done
}

case "$1" in
	purge)
		if [ -e  /etc/acpi/powerbtn.sh.org ]; then
			rm -f /etc/acpi/powerbtn.sh /dev/null 2>&1 || true
			mv /etc/acpi/powerbtn.sh.org /etc/acpi/powerbtn.sh > /dev/null 2>&1 || true
		fi
		remove_service
		
		update-rc.d easyvdr-save-ramdisk remove >/dev/null || exit $?
	;;
	remove|failed-upgrade|abort-install|abort-upgrade)
		remove_service
	;;
	upgrade|disappear)
		#no op
	;;
	*)
		echo "postrm called with unknown argument \`$1'" >&2
		exit 1
	;;
esac




# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
