# e/easyvdr-3.99.00/etc/xdg/autostart/easyvdr-x-launcher.desktop
[Desktop Entry]
Type=Application
Name=easyvdr x launcher
Exec=/usr/lib/vdr/easyvdr-x-launcher


